# personal/myservices/files/rc.local
#!/bin/sh
[ -x /bin/rc.local ]
exit 127
